# supervisor.conf
[program:rep_webserver]
command = ~/dev/rep/run_gunicorn.sh
stdout_logfile = ~/dev/rep/logs/gunicorn_stdout.log
stderr_logfile = ~/dev/rep/logs/gunicorn_stderr.log
